# features/sync/etc/systemd/system/signage-sync.timer
[Unit]
Description=Run Signage Sync every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target