# targets/mpv/etc/systemd/system/getty@tty1.service.d/conflict.conf
[Unit]
Conflicts=signage-chromium.service